# scripts/com.kycscan.discover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kycscan.discover</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/neel/Documents/amikyced/scripts/discover-addresses.mjs</string>
    </array>

    <!-- Run daily at 4 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/neel/Documents/amikyced</string>

    <key>StandardOutPath</key>
    <string>/Users/neel/Documents/amikyced/data/discovery-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/neel/Documents/amikyced/data/discovery-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
